# glacier-pinquery.service
[Unit]
Description=Pin query app
After=lipstick.target
Requires=dbus.socket lipstick.service

[Service]
Type=simple
ExecStart=/usr/bin/glacier-pinquery
Restart=always
RestartSec=10

[Install]
WantedBy=user-session.target

